# bin/.pystart/etc/supervisord.d/example.ini
[program:example-celery]
command=/home/rice/projects/example/scripts/start.sh celery -c etc/internal.ini -c 1 -l info
directory=/root/projects/example
user=1000
group=1000
autostart = true
autorestart = true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stopasgroup = true

[program:example-server]
command=/home/rice/projects/example/scripts/start.sh celery -c etc/internal.ini -b 0.0.0.0:10020
directory=/home/rice/projects/example
user=1000
group=1000
autostart = true
autorestart = true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stopasgroup = true