# bin/.pystart/etc/supervisord.d/example.ini
[program:example]
command=/root/projects/example/scripts/start.sh celery -c etc/internal.ini -c 1 -l info
directory=/root/projects
user=1000
group=1000
autostart = true
autorestart = true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stopasgroup = true
